# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=uwsgi --http 0.0.0.0:80 --module main:app --processes 1 --threads 4 --buffer-size 65535

[program:rsyslogd]
command=rsyslogd -n
